# build/elasticsearch/etc/services.d/consul/run
#!/usr/bin/with-contenv sh

if [ "${DISABLE_CONSUL}" = "true" ]; then
  exit 0
fi

if [ -n "${EC2_IP}" ]; then
  CONSUL_SERF_ADDRESS="${EC2_IP}:${CONSUL_SERF_PORT:-8301}"
  exec consul \
       --retry-join ${CONSUL_SERF_ADDRESS}
fi